# config/aigg-api.service
[Unit]
Description=AIGG Insights API
After=network.target postgresql.service

[Service]
User=cosmos
Group=cosmos
WorkingDirectory=/home/cosmos/aigg-insights
Environment="PATH=/home/cosmos/aigg-insights/venv/bin"
ExecStart=/home/cosmos/aigg-insights/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
